# roles/prometheus/files/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
User=tier3
ExecStart=/usr/bin/node_exporter

[Install]
WantedBy=default.target
